# ansible/roles/uv_mkosi_project_baseline/templates/boot-payload.service.j2
[Unit]
Description={{ project_name }} Python boot payload
After=basic.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/{{ project_name }}
StandardOutput=journal+console
StandardError=journal+console
RemainAfterExit=yes

[Install]
WantedBy=default.target
